# rootfs/etc/services.d/php-fpm/run
#!/bin/execlineb -P

/usr/local/sbin/php-fpm -F